# docker/olab4/supervisor/supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock

[supervisord]
nodaemon=true

#[program:sshd]
#command=/usr/sbin/sshd -D

[program:mariadb]
command=/usr/bin/mariadb-start.sh

[program:httpd]
command=/usr/sbin/httpd -DFOREGROUND

[group:allservices]
#programs=sshd,mariadb,httpd
programs=mariadb,httpd

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock         ; use a unix:// URL  for a unix socket
